# install/picamera.service
[Unit]
Description=Foo

[Service]
ExecStart=/usr/local/picamera
EnvironmentFile=/etc/picamera/env

[Install]
WantedBy=multi-user.target
