# install/pilapse.service
[Unit]
Description=Pilapse

[Service]
ExecStart=/usr/local/bin/pilapse
EnvironmentFile=/etc/pilapse/env

[Install]
WantedBy=multi-user.target
